# ubuntu-jdk-17/supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=0.0.0.0:9001
username=admin
password=123456

[program:sshd]
command=/usr/sbin/sshd -D

